# qr3_91.conf
[program:qr3_91]
command=python3 /home/curator3/qr3_91/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator3/qr3_91/test.err.log
stdout_logfile=/home/curator3/qr3_91/test.err.log
